# third_party/supervisor/paper_reader.ini
[program:paper_reader]
directory=/home/zyh/paper_reader/
command=/usr/bin/uwsgi --ini /home/zyh/paper_reader/third_party/paper_uwsgi.ini
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
